# deploy/bot/pdg-health.timer
[Unit]
Description=Т»Ј 10 тѕєжњЪУиЉСИђТгА PrivDNS Gateway тЂЦт║иУЄфТБђ

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
